# link.list
# Source:Destination
${DOTFILES}/.zsrc.d/.zshrc:${HOME}/.zshrc
